# board/raspberrypicm4io-64/fs/etc/udev/rules.d/90-nk-macsec.rules
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="20a0", ATTR{idProduct}=="4108", RUN+="/opt/nk-macsec/rekey.sh"

